# webconverger/config/chroot_local-includes/etc/network/if-up.d/ping
#!/bin/sh
wget -qO- --timeout=5 --post-data="M=$(md5sum /sys/class/net/eth0/address |
awk '{print $1}')&V=$(dpkg -s webconverger-base |
awk '/^Version: / { print $2 }')" http://ping.webconverger.org
